# vinil.service
[Unit]
Description=Vinil NFC Player
Documentation=https://github.com/wguilherme/vinil
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=wguilherme
WorkingDirectory=/home/wguilherme/projects/raspberry-nfc
ExecStart=/usr/bin/python3 /home/wguilherme/projects/raspberry-nfc/boot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
